# virtualcollaboard/pass-reset-check.cron
* * * * * padmin /home/padmin/resetjob.sh